# com.stealth-sync.daemon.plist
# stealth-sync Daemon Configuration
# This file defines the launchd plist for macOS automatic startup

<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stealth-sync.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- Python interpreter path -->
        <string>/usr/bin/python3</string>
        <!-- stealth-sync daemon script -->
        <string>/tmp/stealth-sync-bootstrap/src/stealth_sync/daemon.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/stealth-sync-daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/stealth-sync-daemon.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <!-- OpenCode DB path (default: ~/.local/share/opencode/opencode.db) -->
        <key>OPENCODE_DB_PATH</key>
        <string>~/.local/share/opencode/opencode.db</string>
        
        <!-- NVIDIA NIM API configuration -->
        <key>NVIDIA_API_KEY</key>
        <string>nvapi-...</string>
        <key>NVIDIA_BASE_URL</key>
        <string>https://integrate.api.nvidia.com/v1</string>
        
        <!-- Polling configuration -->
        <key>POLL_INTERVAL</key>
        <string>10</string>
        
        <!-- Output configuration -->
        <key>OUTPUT_DIR</key>
        <string>docs/opencode-sessions/</string>
        <key>LOGBOOK_PATH</key>
        <string>logbook.stealth.yaml</string>
    </dict>
</dict>
</plist>
